# root/etc/services.d/haste-server/run
#!/usr/bin/with-contenv bash

cd /app/haste-server || exit

exec \
	s6-setuidgid abc /usr/bin/npm start
